# icalsync.timer
[Unit]
Description=Automatically update calendars with icalsync

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
